# {{cookiecutter.repo_name}}/configurations/ansible/roles/supervisor/templates/gunicorn.conf.j2
[program:gunicorn]
command = gunicorn -c gunicorn.conf.py wsgi:application
directory = {{ project_path }}/{{ project_name }}
